# voltlegal.service
[Unit]
Description=VoltLegal Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/voltlegal
ExecStart=/home/ubuntu/voltlegal/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
